# .config/systemd/user/battery-status.service
[Unit]
Description=Charger notify service

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/home/%u/.local/bin/battery-status

[Install]
WantedBy=default.target
