# .config/systemd/user/battery-warning.timer
[Unit]
Description=Battery Low Warning Timer
Requires=battery-warning.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
