# .config/systemd/user/ianny.service
[Unit]
Description=Ianny Service

[Service]
Type=simple
ExecStart=/usr/bin/ianny

[Install]
WantedBy=default.target
