# .config/systemd/user/mic-led-monitor.service
[Unit]
Description=Microphone LED Monitor
# Start as early as possible in the graphical session
After=graphical-session.target
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 0.1
ExecStart=%h/.local/bin/mic-led-monitor
Restart=on-failure
RestartSec=0.5
TimeoutStopSec=5

[Install]
WantedBy=default.target
