# .config/systemd/user/udiskie-automount.service
[Unit]
Description=Udiskie Automount Service

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
